# src/rootfs-overlay/etc/inittab
# 最小 inittab
initdefault:user
sysinit:/etc/init.d/rcS
respawn:user:@terminal
respawn:rescue:/usr/bin/eshell
