# src/rootfs-overlay/etc/init.d/sshd
### BEGIN INIT INFO
# Provides: sshd
# Required-Start:
# Default-Start: user server-headless
### END INIT INFO
. /etc/init.d/rc.common
echo AUDIT sshd_service $1
/usr/bin/start-stop-daemon --service sshd --action "$1" --exec /usr/bin/sshd --pidfile /var/run/sshd.pid --stdout /var/log/sshd.log --stderr /var/log/sshd.log --require /usr/bin/sshd --require /etc/sodex-admin.conf
